# ft_sevices/services/wordpress/sources/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'

[program:phpmyAdmin]
command=/usr/sbin/php-fpm7 --nodaemonize
