# ft_sevices/services/grafana/sources/supervisord.conf
[supervisord]
nodaemon=true

[program:grafana]
command= /grafana/bin/grafana-server
