# ft_sevices/services/inflexdb/sources/supervisord.conf
[supervisord]
nodaemon=true

[program:influxdb]
command= /usr/sbin/influxd -config /etc/influxdb.conf 
