# ft_sevices/services/srcs/ftp/sources/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid

[program:vsftpd]
command=/usr/sbin/vsftpd /etc/vsftpd/vsftpd.conf

[program:telegraf]
command=telegraf -config /etc/telegraf.conf

[eventlistener:processes]
command= /tmp/stop-supervisor.sh
events=PROCESS_STATE_STOPPED, PROCESS_STATE_EXITED, PROCESS_STATE_FATAL